# rasberry_pi/drone-ble.service
[Unit]
Description=Drone BLE Server
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/drone/drone_ble_project/rasberry_pi
ExecStart=/usr/bin/python3 /home/drone/drone_ble_project/rasberry_pi/drone_ble_server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Bluetooth権限のため
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target